# salt/files/etc/profile.d/slurm.sh
S_PATH=/share/apps/slurm/slurm-19.05.5
PATH=$PATH:$S_PATH/bin:$S_PATH/sbin
