# salt/files/etc/systemd/system/fix-eth1-mac-address.service
[Unit]
Before=matlab-server.service

[Service]
ExecStart=/usr/local/bin/fix-eth1-mac-address

[Install]
WantedBy=multi-user.target
