# salt/files/etc/systemd/system/matlab.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/fix-eth1-mac-address

[Install]
WantedBy=matlab-server.target
